# modules/base/files/etc/xinetd.d/talk
service ntalk
{

   disable        = no
   socket_type    = dgram
   wait           = yes
   user           = nobody
   group          = tty
   server         = /usr/sbin/in.ntalkd
   bind           = 127.0.0.1
   only_from      = 127.0.0.1
}
